# deploy/kairoco-backend.service
[Unit]
Description=Kairo Co FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/kairoco-backend
EnvironmentFile=/var/www/kairoco-backend/.env
ExecStart=/var/www/kairoco-backend/.venv/bin/gunicorn app.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --timeout 60
Restart=on-failure

[Install]
WantedBy=multi-user.target
